# launcher.service
[Unit]
Description=My Launcher Service

[Service]
Type=simple
ExecStart=/usr/bin/display.sh

[Install]
WantedBy=multi-user.target
